# app-misc/cortex/files/cortex.service
[Unit]
Description=Cortex - Vault Knowledge Agent (MCP Server)
After=network.target

[Service]
Type=simple
User=cortex
Group=cortex
WorkingDirectory=/var/lib/cortex
EnvironmentFile=-/etc/cortex/cortex.env
ExecStart=/usr/bin/cortex heartbeat --no-bb-poll --format json
Restart=always
RestartSec=300

[Install]
WantedBy=multi-user.target
