# app-misc/cortex/files/cortex.service
[Unit]
Description=Cortex - Vault Knowledge Agent (MCP Server)
After=network.target

[Service]
Type=simple
User=cortex
Group=cortex
WorkingDirectory=/var/lib/cortex
EnvironmentFile=-/etc/cortex/cortex.env
ExecStart=/usr/bin/cortex index --heartbeat # Example: run heartbeat on start
Restart=always
RestartSec=300

[Install]
WantedBy=multi-user.target
