# app-admin/obsidian-manager/files/obsidian-headless@.service
[Unit]
Description=Obsidian Headless Instance (KasmVNC) for vault %i
After=network.target

[Service]
User=obsidian
Group=obsidian
WorkingDirectory=/var/lib/obsidian/%i
Environment=HOME=/var/lib/obsidian/%i
# Use the vault-specific KasmVNC config
ExecStart=/usr/bin/kasmvncserver --config /etc/obsidian/vaults/%i.vnc.yaml --select-de /usr/bin/obsidian --no-sandbox
ExecStop=/usr/bin/kasmvncserver -kill
Restart=always

[Install]
WantedBy=multi-user.target
