# app-misc/qdrant/files/qdrant.service
[Unit]
Description=Qdrant Vector Search Engine
Documentation=https://qdrant.tech/documentation/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/qdrant --config-path /etc/qdrant/config.yaml
WorkingDirectory=/var/lib/qdrant
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/qdrant
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
