# app-misc/hermes/files/hermes-agent.service
[Unit]
Description=Hermes Agent Gateway Daemon
Documentation=https://github.com/NousResearch/hermes-agent
After=network.target

[Service]
Type=simple
User=hermes-webui
Group=hermes-webui
WorkingDirectory=/var/lib/hermes-webui
Environment=HOME=/var/lib/hermes-webui
Environment=XDG_CONFIG_HOME=/var/lib/hermes-webui/.config
ExecStart=/usr/bin/hermes gateway run --accept-hooks
Restart=always
RestartSec=10

# Security sandbox
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
