# app-misc/hermes-webui/files/hermes-webui.service
[Unit]
Description=Hermes Web UI
After=network.target

[Service]
Type=simple
User=hermes-webui
Group=hermes-webui
WorkingDirectory=/var/lib/hermes-webui
EnvironmentFile=-/etc/hermes/webui.conf
ExecStart=/usr/bin/hermes-webui
Restart=always
RestartSec=5

# Security sandbox
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
