# app-admin/mcp-mesh/files/mcp-mesh.service
[Unit]
Description=Model Context Protocol Mesh Gateway
After=network.target

[Service]
Type=simple
User=mcp
Group=mcp
EnvironmentFile=-/etc/conf.d/mcp-mesh
ExecStart=/usr/bin/mcp-mesh serve --port ${PORT} --audit /var/lib/mcp-mesh/audit.jsonl
Restart=on-failure

# --- Hardening (added 2026-06-25, OpenSpec harden-mcp-mesh-service) ---
# Privilege & execution
NoNewPrivileges=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
UMask=0077

# Filesystem
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes

# Kernel & cgroups
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Sockets & namespaces
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
PrivateUsers=yes

# System call filtering
SystemCallFilter=~@system-service
CapabilityBoundingSet=~CAP_SYSLOG

# Managed directories (see mcp-mesh.confd for audit-path interaction)
LogsDirectory=mcp-mesh
ConfigurationDirectory=mcp-mesh
StateDirectory=mcp-mesh

[Install]
WantedBy=multi-user.target
